# pi_py/startup_utils/bpm.service
[Unit]
Description=Start and Update BPM Python Script
Requires=network.target
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/pi/bpm_run/GoogleIoTCoreApp/pi_py/
User=pi
Group=pi
ExecStartPre=/bin/sleep 15
ExecStartPre=/usr/bin/git pull
ExecStartPre=/home/pi/bpm_run/GoogleIoTCoreApp/pi_py/startup_utils/run_pip.sh
ExecStart=/usr/bin/python -u bpm_writer.py
TimeoutStopSec=120


[Install]
WantedBy=multi-user.target